# ansible/templates/pulseaudio.service.j2
[Unit]
Description=PulseAudio Sound Server

[Service]
ExecStart=/usr/bin/pulseaudio --daemonize=yes
Restart=always

[Install]
WantedBy=default.target